# com.neurolearn.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.neurolearn.proxy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/zyc/NeuroLearn/claude_proxy.py</string>
  </array>

  <!-- 登入後自動啟動 -->
  <key>RunAtLoad</key>
  <true/>

  <!-- 若程式崩潰自動重啟 -->
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/neurolearn_proxy.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/neurolearn_proxy.log</string>
</dict>
</plist>
